# infra/templates/reth.service.j2
[Unit]
Description=Reth Ethereum Client
After=network.target
Wants=network.target

[Service]
Type=simple
User=reth
Group=reth
ExecStart=/usr/local/bin/reth node \
  --full \
  --datadir /opt/reth/data \
  --port 30303 \
  --ipcpath /opt/reth/data/reth.ipc \
  --authrpc.addr 0.0.0.0 \
  --authrpc.port 8551 \
  --authrpc.jwtsecret /opt/reth/jwt.hex
WorkingDirectory=/opt/reth
Restart=always
RestartSec=5
TimeoutStopSec=180
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target